# etc/services/named/run
#!/usr/bin/with-contenv sambashare

named -g
